# firmware/console-pi/image-builder/pi-gen-overlay/stage2/99-toypad/files/toypad.service
[Unit]
Description=LEGO Dimensions Toy Pad USB Gadget
After=local-fs.target

[Service]
Type=simple
ExecStartPre=/opt/toypad/setup_gadget.sh setup
ExecStart=/usr/bin/python3 -u /opt/toypad/toy_daemon.py
ExecStopPost=/opt/toypad/setup_gadget.sh teardown
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
